# usbauth/data/20-usbauth.rules
SUBSYSTEM=="usb", ACTION=="add", RUN+="/usr/sbin/usbauth udev-add"

